# web_release.crontab
JAVA_HOME=/software/jdk
LSB_DEFAULTPROJECT=ensembl-core

1 0 * * * $HOME/ensj-healthcheck/run-healthcheck-web-parallel.sh > $HOME/ensj-healthcheck/web.out 2>&1
