# web_release.crontab
JAVA_HOME=/software/jdk1.6.0_01
LSB_DEFAULTPROJECT=ensembl-core

1 0 * * * $HOME/ensj-healthcheck/run-healthcheck-web-parallel.sh > $HOME/ensj-healthcheck/web.out 2>&1

1 0 * * * $HOME/production_dump/dump_production.sh > $HOME/production_dump/cron.log 2>&1
