# web_release.crontab
1 0 * * * /nfs/acari/ensembl/ensj-healthcheck/run-healthcheck-web.sh > /nfs/acari/ensembl/ensj-healthcheck/web.out
