# web_release.crontab
1 0 * * * /nfs/acari/gp1/work/ensj-healthcheck/run-healthcheck-web.sh -debug > /nfs/acari/gp1/work/ensj-healthcheck/web.out
